# build/debian/usr/lib/systemd/system/dbcalm-api.service
[Unit]
Description=Dbcalm
After=network.target

[Service]
User=dbcalm
Group=dbcalm
WorkingDirectory=/usr/bin/
ExecStart=/usr/bin/dbcalm server
Restart=always
Environment="PATH=/usr/bin:/bin"
RuntimeDirectory=dbcalm
RuntimeDirectoryMode=2774

[Install]
WantedBy=multi-user.target